# vagrant/modules/hannibal/templates/etc/profile.d/hannibal.sh
export HANNIBAL_COMPACTIONS_LOGFILE_URL_PATTERN="http://%hostname%:%infoport%/logs/hbase-hbase-master-%hostname%.log"
export HANNIBAL_HBASE_VERSION="<%= hbaseVersion %>"