# vagrant/modules/hannibal/templates/etc/profile.d/hannibal.sh
export HANNIBAL_LOGFILE_PATH_PATTERN="(?i)\"/logs/(.*master.*[.].*)\""
export HANNIBAL_HBASE_VERSION="<%= hbaseVersion %>"
